# nautilus-lab/deploy/systemd/ntlab-forward.timer
[Unit]
Description=Периодический Nautilus-native форвард-снимок (каждые 5 мин)
[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
[Install]
WantedBy=timers.target
